# .docker/crontab
* * * * * /var/www/mykino.top/current/bin/console swiftmailer:spool:send
0 3 * * * /var/www/mykino.top/current/bin/console app:run-release-date-queue
