# .docker/crontab
* * * * * /var/www/mykino.top/current/bin/console swiftmailer:spool:send